# p2pool/rootfs/etc/s6-overlay/s6-rc.d/status_commande/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: P2pool
# Runs status_commande
# ==============================================================================

COMMANDE=status

#   Set variables according to the parameters of sidechain
if bashio::config.true 'p2pool_mini_sidechain'; then
    SIDE_CHAIN_ID=mini
else
    SIDE_CHAIN_ID=default
fi



# --- Wait for p2pool is running ---
while ! pidof -q p2pool; do
    sleep 5
    bashio::log.error "P2POOL IS NOT RUNNING YET"
done


# --- Wait for p2pool sync ---
while (($(jq '.pool.hashrate' /tmp/stats_mod) < 2000000)); do
    bashio::log.error "Wait for p2pool sync"
    sleep 5
done
bashio::log.info "P2pool is now sync with $SIDE_CHAIN_ID SideChain."


#   Set number of second to wait
if ! bashio::config.is_empty 'p2pool_print_status'; then
    DODO=$(( $(bashio::config 'p2pool_print_status') * 60 ))
else
    DODO=900
fi


# --- Execute status command ---
while pidof -q p2pool; do
    sleep $DODO
    echo $COMMANDE
done
