# p2pool/rootfs/etc/s6-overlay/s6-rc.d/p2pool/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: P2pool
# Runs P2pool
# ==============================================================================

declare -a conf_props


# Set monero wallet address
bashio::log.notice "Adding in monero wallet config options..."
conf_props+=(--wallet $(bashio::config 'monero_wallet'))

# Set IP address of your Monero node 
bashio::log.notice "Adding the Monero node address: $(bashio::config 'host')"
conf_props+=(--host $(bashio::config 'host'))


# --- Load USER'S CUSTOM configuration overrides ---
bashio::log.debug "Loading 'p2pool_conf_overrides'..."
for var in $(bashio::config 'p2pool_conf_overrides|keys'); do
    property=$(bashio::config "p2pool_conf_overrides[${var}].property")
    value=$(bashio::config "p2pool_conf_overrides[${var}].value")
    if [[ ${property} =~ secret|password ]]; then
        bashio::log.info "Setting ${property} to ******"
    else
        if ! bashio::config.is_empty "p2pool_conf_overrides[${var}].value"; then
        bashio::log.info "Setting ${property} to ${value}"
        else
        bashio::log.info "Setting ${property} to on"
        fi
    fi
    if ! bashio::config.is_empty "p2pool_conf_overrides[${var}].value"; then
    conf_props+=(--${property} ${value//\"/\\\"})
    else
    conf_props+=(--${property})
    fi
done


# Connect to the p2pool-mini sidechain
if bashio::config.true 'p2pool_mini_sidechain'; then
    conf_props+=(--mini)
    bashio::log.info "Configuration set to connect to the p2pool-mini sidechain"
fi


# Allow UPnP port forwarding
if bashio::config.true 'p2pool_UPnP_port_forwarding'; then
    bashio::log.info "Configuration set to allow UPnP port forwarding"
else 
    conf_props+=(--no-upnp)
    bashio::log.info "Configuration set to not allow UPnP port forwarding"
fi


# Allow UPnP stratum port forwarding
if bashio::config.true 'p2pool_UPnP_portstratum_forwarding'; then
    conf_props+=(--upnp-stratum)
    bashio::log.info "Configuration set to allow UPnPstratum port forwarding"
else 
    bashio::log.info "Configuration set to not allow UPnPstratum port forwarding"
fi


# Set log level
case "$(bashio::config 'p2pool_loglevel')" in \
    0)      loglevel='0' ;; \
    1)      loglevel='1' ;; \
    2)      loglevel='2' ;; \
    3)      loglevel='3' ;; \
    4)      loglevel='4' ;; \
    5)      loglevel='5' ;; \
    6)      loglevel='6' ;; \
esac;
conf_props+=(--loglevel $loglevel)
bashio::log.info "P2pool log level set to $loglevel"


# Launching p2pool
bashio::log.info "Starting of P2Pool"
exec /opt/p2pool --data-api /tmp --local-api "${conf_props[@]}"
