# p2pool/rootfs/etc/services.d/status_commande/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: P2pool
# Runs status_commande
# ==============================================================================

# --- Declare of variables and functions ---

#   Set variables according to the parameters of sidechain
if bashio::config.true 'p2pool_mini_sidechain'; then
    side_chain_ID=mini
    json_pool_info=https://mini.p2pool.observer/api/pool_info
    json_shares_windows=https://mini.p2pool.observer/api/shares_in_window/$(bashio::config 'monero_wallet')
    json_network=https://p2pool.io/mini/api/network/stats
else
    side_chain_ID=default
    json_pool_info=https://p2pool.observer/api/pool_info
    json_shares_windows=https://p2pool.observer/api/shares_in_window/$(bashio::config 'monero_wallet')
    json_network=https://p2pool.io/api/network/stats
fi

#   Set number of second to wait
if ! bashio::config.is_empty 'p2pool_print_status'; then
    dodo=$(( $(bashio::config 'p2pool_print_status') * 60 ))
else
    dodo=900
fi

#   Get Json data from URL
function get_json {
  curl --silent --show-error $1 -H "Accept: application/json" > $2
}

#   Convert variable to float
function to_float {
  printf %.3f "$((10**$4 * $(jq $1 $2) / $3))e-$4"
}

#   Get number of uncles from JSON
function calcul_uncles_sidechain {
  oncles_side=0
  for ChaquesEntree in $(jq -c '.' /tmp/get_json/shares_in_windows.json); do
    if [[ $ChaquesEntree =~ parent ]]; then
        ((oncles_side+=1))
    fi
  done
  printf $oncles_side
}

mkdir /tmp/variables
echo 0 > /tmp/variables/oncles.txt
function calcul_uncles_stratum {
  oncles_stratum=$(</tmp/variables/oncles.txt)
  touch /tmp/variables/timestamp.json
  for ChaquesEntree in $(jq -c '.[]' /tmp/get_json/shares_in_windows.json); do
    if [[ $ChaquesEntree =~ parent ]]; then
      echo $ChaquesEntree > /tmp/variables/timestamp.json
      if [[ $(( $(date +%s) - $dodo )) -le $(jq '. | .timestamp' /tmp/variables/timestamp.json) ]]; then
        ((oncles_stratum+=1))
      fi
    fi
  done
  echo $oncles_stratum > /tmp/variables/oncles.txt
  printf $oncles_stratum
}



# --- Wait for p2pool is running ---
while ! pidof -q p2pool; do
    sleep 5
    bashio::log.error "P2POOL IS NOT RUNNING YET"
done



# --- Save p2pool PID ---
pidof p2pool > /tmp/p2pool.pid

PID=$(</tmp/p2pool.pid)

sleep 5



# --- Wait for p2pool sync ---
while (($(jq '.pool.hashrate' /tmp/stats_mod) < 2000000)); do
    bashio::log.error "Wait for p2pool sync"
    sleep 5
done
bashio::log.info "P2pool is now sync with $side_chain_ID SideChain."



# --- Execute commande status command ---
while pidof -q p2pool; do
    pidof p2pool > /tmp/p2pool.pid
    PID=$(</tmp/p2pool.pid)
    # --- SideChain status ---
    mkdir /tmp/get_json
    get_json $json_pool_info /tmp/get_json/pool_info.json
    get_json $json_shares_windows /tmp/get_json/shares_in_windows.json
    get_json $json_network /tmp/get_json/network_stats.json
    bashio::log.notice "SideChain status"
    echo "Main chain height               = $(jq '.mainchain.height' /tmp/get_json/pool_info.json)"
    echo "Main chain hashrate             = $(to_float '.difficulty' /tmp/get_json/network_stats.json $(($(jq '.mainchain.block_time' /tmp/get_json/pool_info.json) * 1000000000)) 6) GH/s"
    echo "Side chain ID                   = $side_chain_ID"
    echo "Side chain height               = $(jq '.sidechain.height' /tmp/get_json/pool_info.json)"
    echo "Side chain hashrate             = $(to_float '.pool.hashrate' /tmp/stats_mod 1000000 6) MH/s"
    echo "PPLNS window                    = $(jq '.sidechain.window.blocks' /tmp/get_json/pool_info.json) blocks (+$(jq '.sidechain.window.uncles' /tmp/get_json/pool_info.json) uncles)"
    echo "Shares in current PPNLS windows = $(( $(jq '. | length' /tmp/get_json/shares_in_windows.json) - $(calcul_uncles_sidechain) )) blocks (+$(calcul_uncles_sidechain) uncles)"
    # --- StratumServer status ---
    cpt=0
    while ((cpt<7200)); do    
        bashio::log.info "StratumServer status"
        if test -f "/tmp/local/stats"; then
            echo "Hashrate (15m est) = $(to_float '.hashrate_15m' /tmp/local/stats 1000 3) KH/s"
            echo "Hashrate (1h est)  = $(to_float '.hashrate_1h' /tmp/local/stats 1000 3) KH/s"
            echo "Hashrate (24h est) = $(to_float '.hashrate_24h' /tmp/local/stats 1000 3) KH/s"
            echo "Total hashes       = $(jq '.total_hashes' /tmp/local/stats)"
            calcul_uncles_stratum
            echo "Shares found       = $(( $(jq '.shares_found' /tmp/local/stats) - $(</tmp/variables/oncles.txt) )) (+$(</tmp/variables/oncles.txt) uncles)"
            echo "Average effort     = $(jq '.average_effort' /tmp/local/stats)%"
            echo "Current effort     = $(jq '.current_effort' /tmp/local/stats)%"
            echo "Connections        = $(jq '.connections' /tmp/local/stats) ($(jq '.incoming_connections' /tmp/local/stats) incoming)"
            rm /tmp/local/stats
        else
            echo "No active miner connexion"
            echo "Hashrate (15m est) = 0 KH/s"
            echo "Hashrate (1h est)  = 0 KH/s"
            echo "Hashrate (24h est) = 0 KH/s"
            echo "Total hashes       = 0"
            echo "Shares found       = 0"
            echo "Average effort     = 0.000%"
            echo "Current effort     = 0.000%"
            echo "Connections        = 0 (0 incoming)"
        fi
        sleep $dodo
        ((cpt+=$dodo))
    done
done
